# service_startup/button-control.service
[Unit]
Description=KY-004 Button Control (Start Service & Shutdown)
After=multi-user.target

[Service]
Type=simple
User=your-user
WorkingDirectory=/home/your-user/dogo-cam
ExecStart=/usr/bin/python3 /home/your-user/dogo-cam/ky004-control.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
