# service_startup/cloudflared-tunnel.service
[Unit]
Description=Cloudflare Tunnel for Dog Stream
After=network.target

[Service]
User=your-user
EnvironmentFile=-/home/your-user/dogo-cam/.env
Environment=CLOUDFLARED_CONFIG=/home/your-user/.cloudflared/config.yml
Environment=CLOUDFLARED_TUNNEL=dog-stream-tunnel
ExecStart=/usr/local/bin/cloudflared tunnel --config ${CLOUDFLARED_CONFIG} run ${CLOUDFLARED_TUNNEL}
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
