# service_startup/dog-stream-flask.service
[Unit]
Description=Dog Stream Flask App
After=network.target

[Service]
User=your-user
WorkingDirectory=/home/your-user/dogo-cam
EnvironmentFile=/home/your-user/dogo-cam/.env
ExecStart=/home/your-user/.local/bin/uv run gunicorn --worker-class gthread --workers 1 --threads 4 --bind 0.0.0.0:5000 dogcam_stream:app
Restart=always
RestartSec=10
LimitNOFILE=4096
OOMScoreAdjust=-1000

[Install]
# Service is started manually by button-control.service
# Not auto-started on boot to save power
# WantedBy=multi-user.target
